# temperature2mqtt.service
[Unit]
Description=Temperature to MQTT gateway
After=network.target

[Service]
Type=simple
ExecStart=/usr/local/bin/temperature2mqtt.sh
StandardOutput=syslog
StandardError=syslog

[Install]
WantedBy=multi-user.target
